# logrotate.cron
#!/bin/bash

/usr/sbin/logrotate  --state=/logrotate-status/logrotate.status /logrotate-status/logrotate.conf 
